# jupyterhub.service
[Unit]
Description=JupyterHub Multi-User Server
Documentation=https://jupyterhub.readthedocs.io/
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/workspace/jupyterhub
ExecStart=/openhands/poetry/openhands-ai-5O4_aCHf-py3.12/bin/jupyterhub --config=jupyterhub_config_custom.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jupyterhub

# Environment variables
Environment=PATH=/openhands/poetry/openhands-ai-5O4_aCHf-py3.12/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/workspace/jupyterhub

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/workspace/jupyterhub

[Install]
WantedBy=multi-user.target